# docker/crontab
0 * * * * root curl -sf http://localhost/run -H "Authorization: Bearer $CRON_SECRET" >> /proc/1/fd/1 2>&1
