# docker/crontab
0 * * * * root curl -sf http://localhost/run -H "Authorization: Bearer $CRON_SECRET" >> /var/log/kinsync-cron.log 2>&1
